# scripts/launchagent/work.jowork.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS LaunchAgent template for Jowork gateway.
  Copy to ~/Library/LaunchAgents/work.jowork.plist
  then run: launchctl load ~/Library/LaunchAgents/work.jowork.plist

  Replace {{DATA_DIR}} and {{APP_DIR}} with actual paths before loading.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>work.jowork</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>{{APP_DIR}}/dist/cluster.js</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>DATA_DIR</key>
    <string>{{DATA_DIR}}</string>
    <key>PORT</key>
    <string>18800</string>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>

  <!-- Raise open-file limit to handle many concurrent connections -->
  <key>SoftResourceLimits</key>
  <dict>
    <key>NumberOfFiles</key>
    <integer>8192</integer>
  </dict>
  <key>HardResourceLimits</key>
  <dict>
    <key>NumberOfFiles</key>
    <integer>8192</integer>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>{{DATA_DIR}}/logs/jowork.log</string>

  <key>StandardErrorPath</key>
  <string>{{DATA_DIR}}/logs/jowork-error.log</string>

  <key>WorkingDirectory</key>
  <string>{{APP_DIR}}</string>
</dict>
</plist>
